# housing-deb/src/deb/control/prerm
#!/bin/sh
systemctl stop    housing || true
systemctl disable housing || true
